# borg.service
[Unit]
Description=Start backup script

[Service]
Type=oneshot
ExecStart=/bin/sh /vagrant/backup.sh